# apache.sv.conf
[program:apache]
command=/usr/sbin/apache2ctl start
autorestart=true
startsecs=5
autostart=true