# ssh.sv.conf
[program:ssh]
command=/usr/sbin/sshd
autorestart=true
startsecs=5
autostart=true